# eda-filewatch@.service
[Unit]
Description=EDA File Watch Monitor for AAP (%i)
Documentation=https://github.com/user/eda-filewatch-systemd
After=network.target
Wants=network.target

[Service]
Type=simple
User=eda-filewatch
Group=eda-filewatch
WorkingDirectory=/opt/eda-filewatch
ExecStart=/opt/eda-filewatch/filewatch-monitor.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5

# Environment configuration
Environment=CONFIG_FILE=/etc/eda-filewatch/%i.conf
Environment=LOG_LEVEL=INFO

# Security settings
NoNewPrivileges=true
PrivateTmp=false
PrivateDevices=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/var/log/eda-filewatch /tmp
ReadOnlyPaths=/etc/eda-filewatch /home /var /opt /usr

# Resource limits
LimitNOFILE=1024
LimitNPROC=32

[Install]
WantedBy=multi-user.target 